# rmilter.socket
[Unit]
Description=Another sendmail milter for different mail checks

[Install]
WantedBy=sockets.target
