# debian/postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ]; then
   update-rc.d rmilter remove >/dev/null
fi

if [ -d /run/systemd/system ]; then
    systemctl --system daemon-reload >/dev/null || true
fi

#DEBHELPER#
