# debian/prerm
#!/bin/sh
set -e

if [ -x /etc/init.d/rmilter ] && ! [ -d /run/systemd/system ]; then
    invoke-rc.d rmilter stop || exit $?
fi

#DEBHELPER#
