# debian/rmilter.service
[Unit]
Description=Another sendmail milter for different mail checks

[Service]
ExecStart=/usr/sbin/rmilter -c /etc/rmilter.conf -n
User=_rmilter
