# debian/rmilter.upstart
description "rmilter - sendmail filter"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

exec start-stop-daemon -S -q -c _rmilter -x /usr/sbin/rmilter -- -c /etc/rmilter.conf.sysvinit -n

pre-start script
	install -o _rmilter -g adm -m 750 -d /run/rmilter
	rm -f /run/rmilter/rmilter.sock
	chmod g+s /run/rmilter
end script
